# roles/patwinlink/templates/patwinlink.target.j2
[Unit]
Description=Pat Winlink service group
Requires=rigctld.service ardopc.service pat-ardopc.service
After=rigctld.service ardopc.service pat-ardopc.service
